# custom/openrc_fandaemon
#!/sbin/openrc-run

description="HDD temp based Fan Daemon"

depend()
{
	after miconstartup
	use modules devfs
}

start()
{
  . /etc/buffalo_type
  start-stop-daemon --start --quiet --background --make-pidfile --pidfile /var/run/fand.pid --exec /usr/bin/fandaemon_wrapper.sh "$fan_type"
}

stop()
{
  start-stop-daemon --stop --quiet --pidfile /var/run/fand.pid
}
